# scripts/task_check_threadtracking_health.xml
<?xml version="1.0" encoding="UTF-8"?>
<Task version="1.4" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
  <RegistrationInfo>
    <Description>Checks for missing ThreadTracking rows and logs a health report. Uses Django command: manage.py check_threadtracking_health</Description>
    <Author>GmailJobTracker</Author>
  </RegistrationInfo>
  <Triggers>
    <CalendarTrigger>
      <!-- Adjust this start time after import if desired -->
      <StartBoundary>2025-11-10T09:00:00</StartBoundary>
      <Enabled>true</Enabled>
      <ScheduleByDay>
        <DaysInterval>1</DaysInterval>
      </ScheduleByDay>
    </CalendarTrigger>
  </Triggers>
  <Principals>
    <Principal id="Author">
      <LogonType>InteractiveToken</LogonType>
      <RunLevel>LeastPrivilege</RunLevel>
    </Principal>
  </Principals>
  <Settings>
    <MultipleInstancesPolicy>IgnoreNew</MultipleInstancesPolicy>
    <DisallowStartIfOnBatteries>false</DisallowStartIfOnBatteries>
    <StopIfGoingOnBatteries>false</StopIfGoingOnBatteries>
    <AllowHardTerminate>true</AllowHardTerminate>
    <StartWhenAvailable>true</StartWhenAvailable>
    <RunOnlyIfNetworkAvailable>false</RunOnlyIfNetworkAvailable>
    <IdleSettings>
      <StopOnIdleEnd>false</StopOnIdleEnd>
      <RestartOnIdle>false</RestartOnIdle>
    </IdleSettings>
    <AllowStartOnDemand>true</AllowStartOnDemand>
    <Enabled>true</Enabled>
    <Hidden>false</Hidden>
    <RunOnlyIfIdle>false</RunOnlyIfIdle>
    <WakeToRun>false</WakeToRun>
    <ExecutionTimeLimit>PT30M</ExecutionTimeLimit>
    <Priority>7</Priority>
  </Settings>
  <Actions Context="Author">
    <Exec>
      <Command>C:\Users\kaver\code\GmailJobTracker\.venv\Scripts\python.exe</Command>
      <Arguments>manage.py check_threadtracking_health --days 7 --autofix</Arguments>
      <WorkingDirectory>C:\Users\kaver\code\GmailJobTracker</WorkingDirectory>
    </Exec>
  </Actions>
  <!-- Logs are written by the command itself to logs\\threadtracking_health.log -->
</Task>
